# labwc-config/systemd/labwc-theme-watcher.service
[Unit]
Description=Watch Labwc theme files and reconfigure on changes
After=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'inotifywait -m -e modify,create ~/.local/share/themes/Vermello/openbox-3/themerc ~/.config/labwc/labwc.xml 2>/dev/null | while read; do labwc --reconfigure; done'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
